# upstart/harold-startup.conf
description "start all described instances of harold"

task
start on runlevel [2345]

script
    for ini in /opt/harold/etc/*.ini; do
        start --no-wait harold ini=$ini
    done
end script
