# upstart/harold-startup.conf
description "start all described instances of harold"

task
start on runlevel [2345]

script
    initctl emit harold-action ACTION=start
end script
