# upstart/harold-prep.conf
description "ensure /var/run/harold is present and writable by harold"

task
start on starting JOB=harold

script
	mkdir -p /var/run/harold
	chown harold:harold /var/run/harold
end script
